# templates/ipfs.service.j2
[Unit]
Description=IPFS Daemon


[Service]
ExecStart=su - ipfs bash -c "/usr/local/bin/ipfs daemon"


[Install]
WantedBy=multi-user.target
